# files/niri/usr/etc/xdg/autostart/nvidia-settings-load.desktop
[Desktop Entry]
Type=Application
Name=NVIDIA X Server Settings
Comment=Disabled in the Niri image: nvidia-settings --load-config-only exits non-zero under the Wayland/Niri session
Exec=/usr/bin/nvidia-settings --load-config-only
Hidden=true
X-GNOME-Autostart-enabled=false
