# files/niri/usr/lib/systemd/system/systemd-remount-fs.service.d/10-skip-composefs.conf
# Fedora 44 composefs deployments expose / as a read-only overlay.
# systemd-remount-fs cannot reconfigure that mount, so skip the unit there.
[Unit]
ConditionPathExists=!/run/ostree/.private/cfsroot-lower
