# files/system/usr/lib/systemd/system/bluetooth.service.d/10-config-dir-mode.conf
[Service]
# Fedora bluez ships /etc/bluetooth as 0755, while the unit asks systemd for a
# 0555 ConfigurationDirectory. Match the packaged directory mode to avoid a
# noisy boot warning without changing bluetoothd behavior.
ConfigurationDirectoryMode=0755
