# files/system/usr/lib/systemd/user/syncthing.service.d/override.conf
[Service]
# Force the GUI to bind to the default port (8384) instead of the random
# port syncthing picks when generating a fresh config.xml on first launch.
Environment=STGUIADDRESS=127.0.0.1:8384
