# files/system/usr/lib/systemd/user/syncthing.service.d/override.conf
[Unit]
# This user unit is enabled globally, so greeter/system user managers would
# otherwise try to start Syncthing with homes such as /var/lib/greetd.
ConditionUser=!@system
ConditionUser=!greetd
ConditionUser=!sddm

[Service]
# Force the GUI to bind to the default port (8384) instead of the random
# port syncthing picks when generating a fresh config.xml on first launch.
Environment=STGUIADDRESS=127.0.0.1:8384
