# files/niri/usr/lib/systemd/user/darkman.service.d/10-skip-greeter.conf
[Unit]
# darkman is enabled globally for user sessions, but greeter users do not have
# the desktop/session environment it expects. Avoid noisy failed starts before
# the real user session begins.
ConditionUser=!greetd
ConditionUser=!sddm
