# files/niri/usr/lib/systemd/user/darkman.service.d/10-skip-sddm.conf
# darkman is enabled globally for interactive users, but the SDDM greeter
# has no darkman configuration and repeatedly fails the user unit.
[Unit]
ConditionUser=!sddm
