# files/system/usr/lib/systemd/system/thermald.service.d/10-intel-only.conf
[Unit]
# thermald is an Intel thermal daemon. On AMD systems it exits successfully but
# logs "Unsupported cpu model or platform" on every boot.
ConditionPathExists=/sys/devices/system/cpu/intel_pstate
